# agent/service/systemd/breeze-agent.service
[Unit]
Description=Breeze RMM Agent
Documentation=https://github.com/breeze-rmm/breeze
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# breeze-unit-version: 2
Type=simple
ExecStart=/usr/local/bin/breeze-agent start
WorkingDirectory=/etc/breeze
Restart=on-failure
# 30s cooldown spreads respawn across a fleet that crashes simultaneously
# (e.g. correlated network blip). Combined with StartLimitBurst=5 over
# StartLimitIntervalSec=60, a misbehaving host backs off entirely instead
# of stampeding the API.
RestartSec=30

# Cap total stop time so a hung HTTP flush during OS shutdown (network
# going down) doesn't block system power-off for the 90s systemd default.
# KillMode=mixed sends SIGTERM to the main process, then SIGKILL to the
# whole cgroup after TimeoutStopSec.
TimeoutStopSec=15
KillMode=mixed

# INTENTIONALLY UNSANDBOXED. The remote terminal and remote script execution
# features spawn child processes that must behave like a root SSH session:
#   - package managers drop privileges to unprivileged users (needs CAP_SETUID/SETGID/CHOWN)
#   - admins write under /home, /usr, /etc, and expect a shared /tmp
# systemd sandbox restrictions are INHERITED by those children and silently break
# these operations. Do not re-add them.
# See docs/superpowers/specs/2026-06-09-agent-systemd-sandbox-remote-terminal-design.md

StandardOutput=journal
StandardError=journal
SyslogIdentifier=breeze-agent
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
